# deploy/model-price-backend.service
[Unit]
Description=Model Price Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/model-price/backend
Environment=RELOAD=false
Environment=CORS_ORIGINS=["https://price.pubwhere.cn"]
ExecStart=/root/.local/bin/uv run uvicorn main:app --host 127.0.0.1 --port 8090
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
